# data/sbud.service.in
[Unit]
Description=SBU Daemon

[Service]
Type=dbus
BusName=com.hughski.PowerSBU
ExecStart=@LIBEXECDIR@/sbud

[Install]
WantedBy=multi-user.target
